# bootstrap/templates/deepseek-bridge.service
[Unit]
Description=DeepSeek V4 OpenCode Claude Code Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/leonardo/dev/repos/omni-harness
ExecStart=/home/leonardo/.nvm/versions/node/v24.14.1/bin/node /home/leonardo/dev/repos/omni-harness/index.js --config /home/leonardo/dev/repos/omni-harness/config.json
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
